# deploy/scrapower.service
[Unit]
Description=Scrapower Coordinator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=scrapower
Group=scrapower
WorkingDirectory=/opt/scrapower
EnvironmentFile=/opt/scrapower/env
ExecStart=/opt/scrapower/.venv/bin/python -m scrapower.coordinator.main

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/scrapower/data
PrivateTmp=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
